# cronjob/uber_job.cron
* * * * * /usr/bin/python /root/rahul/grab/Surge-Pricing-Calculations-master/cronjobuber_job_scheduler.py > /home/rahul/uber.log

